# util_scripts/verify_systemd_services/systemd_test_suite_full/broken/broken13-missing-service-section.service
[Unit]
Description=No Service Section
Wants=network-online.target
After=network-online.target
After=NetworkManager-wait-online.service
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
