# util_scripts/verify_systemd_services/systemd_test_suite_full/broken/broken01-missing-install.service
[Unit]
Description=Missing Install Section
Wants=network-online.target
After=network-online.target
After=NetworkManager-wait-online.service
StartLimitIntervalSec=0

[Service]
User=root
ExecStart=/bin/true
Restart=always
RestartSec=5
