# atomp_automation_deployment/util_scripts/verify_systemd_services/systemd_test_suite_full/broken/broken12-wantedby-wrong-target.service
[Unit]
Description=Wrong WantedBy Target
Wants=network-online.target
After=network-online.target
After=NetworkManager-wait-online.service
StartLimitIntervalSec=0

[Service]
User=root
ExecStart=/bin/true
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
