# atomp_automation_deployment/util_scripts/verify_systemd_services/systemd_test_suite_full/broken/broken19-tabs-and-spaces.service
[Unit]
	Description=Tabs and Spaces
	Wants=network-online.target
	After=network-online.target
	After=NetworkManager-wait-online.service
	StartLimitIntervalSec=0

[Service]
	User=root
	ExecStart=/bin/true
	Restart=always
	RestartSec=5

[Install]
	WantedBy=multi-user.target
