# atomp_automation_deployment/util_scripts/verify_systemd_services/systemd_test_suite_full/broken/broken06-missing-startlimit.service
[Unit]
Description=Missing StartLimitInterval
Wants=network-online.target
After=network-online.target
After=NetworkManager-wait-online.service

[Service]
User=root
ExecStart=/bin/true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
