# atomp_automation_deployment/util_scripts/verify_systemd_services/systemd_test_suite_full/broken/broken04-user-in-unit.service
[Unit]
Description=User Misplaced
Wants=network-online.target
After=network-online.target
After=NetworkManager-wait-online.service
StartLimitIntervalSec=0
User=root   # WRONG section

[Service]
ExecStart=/bin/true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
